# packaging/systemd/exwm-vr-compositor.service
[Unit]
Description=EXWM-VR Wayland Compositor
Documentation=https://github.com/Jesssullivan/XoxdWM
PartOf=exwm-vr.target

[Service]
Type=simple
ExecStart=/usr/bin/ewwm-compositor
Restart=on-failure
RestartSec=3
Environment=XDG_CURRENT_DESKTOP=EXWM-VR
Environment=XDG_SESSION_TYPE=wayland

[Install]
WantedBy=exwm-vr.target
