# packaging/systemd/exwm-vr-emacs.service
[Unit]
Description=EXWM-VR Emacs Window Manager
Documentation=https://github.com/Jesssullivan/XoxdWM
After=exwm-vr-compositor.service
BindsTo=exwm-vr-compositor.service
Requires=exwm-vr-compositor.service
PartOf=exwm-vr.target

[Service]
Type=simple
ExecStart=/usr/bin/emacs --fg-daemon=exwm-vr -l ewwm
Environment=QT_QPA_PLATFORM=wayland
Environment=MOZ_ENABLE_WAYLAND=1
Environment=_JAVA_AWT_WM_NONREPARENTING=1
Environment=GDK_BACKEND=wayland,x11
Environment=XCURSOR_THEME=Adwaita
Environment=XCURSOR_SIZE=24

[Install]
WantedBy=exwm-vr.target
