# packaging/systemd/exwm-vr-monado.service
[Unit]
Description=EXWM-VR Monado OpenXR Runtime
Documentation=https://github.com/Jesssullivan/XoxdWM
After=exwm-vr-compositor.service
ConditionPathExists=/etc/xdg/openxr/1/active_runtime.json
PartOf=exwm-vr.target

[Service]
Type=simple
ExecStart=/usr/bin/monado-service
Environment=XR_RUNTIME_JSON=/etc/xdg/openxr/1/active_runtime.json
Restart=on-failure

[Install]
WantedBy=exwm-vr.target
