# packaging/systemd/exwm-vr-brainflow.service
[Unit]
Description=EXWM-VR BrainFlow BCI Daemon
Documentation=https://github.com/Jesssullivan/XoxdWM
After=exwm-vr-compositor.service
ConditionPathExists=/dev/openbci
PartOf=exwm-vr.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/true

[Install]
WantedBy=exwm-vr.target
