# packaging/udev/98-bigscreen-non-desktop.rules
# Bigscreen Beyond non-desktop flag setter
#
# On DRM hotplug, scan DP connectors for Beyond EDID and set non_desktop=1
# so compositors treat it as a VR display.
#
# Note: "BIG" is the binary EDID manufacturer code (not strings-extractable).
# We match on the product name "Beyond" which is a plain ASCII string in EDID.
#
# Workaround until the upstream kernel EDID quirk is merged:
#   EDID_QUIRK('B', 'I', 'G', 0x1234, EDID_QUIRK_NON_DESKTOP)
#
# Note: $$ is required to escape $ in udev RUN+= values.
# Note: non_desktop sysfs attr only exists if kernel/driver supports it.

ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/bin/sh -c 'for c in /sys/class/drm/card*-DP-*/; do if [ -f $$c/edid ] && strings $$c/edid 2>/dev/null | grep -q Beyond; then [ -w $$c/non_desktop ] && echo 1 > $$c/non_desktop 2>/dev/null; fi; done'"
